# systemd/cirrus-beat.service
# Cirrus Beat — Celery Beat Scheduler
# Deploy: sudo cp /var/www/cirrus/systemd/cirrus-beat.service /etc/systemd/system/
#         sudo systemctl daemon-reload && sudo systemctl enable cirrus-beat

[Unit]
Description=Cirrus Celery Beat
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=farizpe
Group=farizpe
WorkingDirectory=/var/www/cirrus
Environment="PATH=/var/www/cirrus/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/var/www/cirrus/.env

ExecStart=/var/www/cirrus/venv/bin/celery \
    -A cirrus beat \
    -l info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler

Restart=on-failure
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/cirrus/logs
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
